# config/systemd/system/sepolia.target
[Unit]
Description=Sepolia local node
Wants=sepolia-reth.service sepolia-lighthouse.service prometheus.service grafana.service
